# prometheus.conffiles
/etc/prometheus/prometheus.conf
/etc/prometheus/hosts
/etc/cron.d/prometheus
